# scripts/systemd/leadgen-nightly.service
[Unit]
Description=Lead Generation Nightly Runner
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/lead_generation-main
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python scripts/run_orchestrator.py --tenant-id 1 --input "run full pipeline"
Restart=on-failure

[Install]
WantedBy=default.target
