# scripts/systemd/leadgen-nightly.service
[Unit]
Description=Lead Generation Nightly Runner
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/lead_generation-main
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python -m scripts.run_nightly
Restart=on-failure

[Install]
WantedBy=default.target

